# docker/supervisord.conf
; GeoPulse all-in-one process tree.
;
; Startup order is by `priority`: Postgres first, then a one-shot migration that
; drops a marker file, then the app services (which block on that marker via
; geopulse-run.sh) and nginx. Every long-running process is restarted on exit.

[supervisord]
nodaemon=true
user=root
pidfile=/run/supervisord.pid
logfile=/dev/null
logfile_maxbytes=0
loglevel=info

[unix_http_server]
file=/run/supervisor.sock

[supervisorctl]
serverurl=unix:///run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; ---- database ----
[program:postgres]
command=docker-entrypoint.sh postgres
priority=10
autostart=true
autorestart=true
stopsignal=INT
stopwaitsecs=40
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

; ---- one-shot: apply SQL migrations, then touch /run/geopulse/migrated ----
[program:migrate]
command=/usr/local/bin/geopulse-migrate.sh
directory=/app
priority=20
autostart=true
autorestart=false
startsecs=0
exitcodes=0
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

; ---- API (localhost only; nginx proxies /api to it) ----
[program:api]
command=/usr/local/bin/geopulse-run.sh python3 -m uvicorn api.main:app --host 127.0.0.1 --port 8000
directory=/app
priority=30
autostart=true
autorestart=true
startsecs=5
startretries=30
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

; ---- ingestion workers + scoring ----
[program:markets-worker]
command=/usr/local/bin/geopulse-run.sh python3 markets-worker/main.py
directory=/app
priority=30
autostart=true
autorestart=true
startsecs=5
startretries=30
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:gdelt-worker]
command=/usr/local/bin/geopulse-run.sh python3 gdelt-worker/main.py
directory=/app
priority=30
autostart=true
autorestart=true
startsecs=5
startretries=30
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:fx-worker]
command=/usr/local/bin/geopulse-run.sh python3 fx-worker/main.py
directory=/app
priority=30
autostart=true
autorestart=true
startsecs=5
startretries=30
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:macro-worker]
command=/usr/local/bin/geopulse-run.sh python3 macro-worker/main.py
directory=/app
priority=30
autostart=true
autorestart=true
startsecs=5
startretries=30
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:weather-worker]
command=/usr/local/bin/geopulse-run.sh python3 weather-worker/main.py
directory=/app
priority=30
autostart=true
autorestart=true
startsecs=5
startretries=30
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:weather-field-worker]
command=/usr/local/bin/geopulse-run.sh python3 weather-field-worker/main.py
directory=/app
priority=30
autostart=true
autorestart=true
startsecs=5
startretries=30
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:storm-worker]
command=/usr/local/bin/geopulse-run.sh python3 storm-worker/main.py
directory=/app
priority=30
autostart=true
autorestart=true
startsecs=5
startretries=30
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:scoring-engine]
command=/usr/local/bin/geopulse-run.sh python3 scoring-engine/main.py
directory=/app
priority=30
autostart=true
autorestart=true
startsecs=5
startretries=30
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

; ---- web (SPA + reverse proxy) ----
[program:nginx]
command=nginx -g "daemon off;"
priority=40
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
